# Source/client/etc/cron.d/snowman-client
#
# Auto-update for the snowman-client.
# 	Default: Every night at 2:00
#
0 2 * * * root /usr/bin/snowmanclient-update reloadconfig
